# tcsd.service
[Unit]
Description=TrouSerS TSS daemon

[Service]
Type=forking
ExecStart=/usr/sbin/tcsd

[Install]
WantedBy=multi-user.target
